# packaging/deb/foldops-agent/prerm
#!/bin/sh
set -e

#DEBHELPER#

case "$1" in
    remove|upgrade|deconfigure)
        if command -v deb-systemd-helper >/dev/null 2>&1; then
            deb-systemd-helper mask foldops-agent.service >/dev/null 2>&1 || true
            deb-systemd-helper stop foldops-agent.service >/dev/null 2>&1 || true
        elif command -v systemctl >/dev/null 2>&1; then
            systemctl stop foldops-agent.service >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
