# systemd/rust/foldops-agent.service
[Unit]
Description=FoldOps FAH farm monitoring agent (Rust)
After=network-online.target fah-client.service
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/foldops/agent.env
ExecStart=/usr/bin/foldops-agent
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
